# etc/skel/.config/autostart/live-usb-storage.desktop
[Desktop Entry]
Name=live-usb-storage
Exec=/usr/local/sbin/live-usb-storage.sh
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
