# etc/skel/.config/autostart/setup-live.desktop
[Desktop Entry]
Name=live-usb-storage
Exec=/usr/local/sbin/setup-live
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
